# dot_config/systemd/user/hyprdynamicmonitors.service
[Unit]
Description=HyprDynamicMonitors - Dynamic monitor configuration for Hyprland
Documentation=https://fiffeek.github.io/hyprdynamicmonitors/
PartOf=graphical-session.target
Requires=graphical-session.target
After=graphical-session.target
After=hyprdynamicmonitors-prepare.service
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=simple
ExecStart=/usr/local/bin/hyprdynamicmonitors run --enable-lid-events
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
